# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:8001]
command = /p.py 0.0.0.0 8001 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8001.log

[program:8002]
command = /p.py 0.0.0.0 8002 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8002.log

[program:8003]
command = /p.py 0.0.0.0 8003 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8003.log

[program:8004]
command = /p.py 0.0.0.0 8004 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8004.log

[program:8005]
command = /p.py 0.0.0.0 8005 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8005.log

[program:8006]
command = /p.py 0.0.0.0 8006 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8006.log

[program:8007]
command = /p.py 0.0.0.0 8007 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8007.log

[program:8008]
command = /p.py 0.0.0.0 8008 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8008.log

[program:8009]
command = /p.py 0.0.0.0 8009 huabei-pool.ethfans.org 13333 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8009.log

[program:8010]
command = /p.py 0.0.0.0 8010 huabei-pool.ethfans.org 13333 0x96d30b101dafb0b70be1c5951da018fe02607357
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8010.log

[program:8011]
command = /p.py 0.0.0.0 8011 huabei-pool.ethfans.org 13333 0xed6d02456b4d0ad5f63e86046d9b4c90c168c023
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8011.log

[program:8012]
command = /p.py 0.0.0.0 8012 huabei-pool.ethfans.org 13333 0xed6d02456b4d0ad5f63e86046d9b4c90c168c023
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8012.log

[program:8013]
command = /p.py 0.0.0.0 8013 huabei-pool.ethfans.org 13333 0xed6d02456b4d0ad5f63e86046d9b4c90c168c023
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8013.log

[program:8014]
command = /p.py 0.0.0.0 8014 huabei-pool.ethfans.org 13333 0xed6d02456b4d0ad5f63e86046d9b4c90c168c023
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8014.log

[program:8015]
command = /p.py 0.0.0.0 8015 huabei-pool.ethfans.org 13333 0xed6d02456b4d0ad5f63e86046d9b4c90c168c023
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8015.log

[program:8016]
command = /p.py 0.0.0.0 8016 huabei-pool.ethfans.org 13333 0xed6d02456b4d0ad5f63e86046d9b4c90c168c023
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8016.log

[program:8017]
command = /p.py 0.0.0.0 8017 pool.sklinux.com 8008 0xe32e3b302c1e30521ccf234de0bb65ec32ded13c
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/supervisor/8017.log

